# packages/deploy-pckg-0.1/debian/deploy-pckg/DEBIAN/postinst
#!/bin/sh
# postinst script for deb-pckgs
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
        [ ! -d /etc/prometheus/ssl ] && mkdir -p /etc/prometheus/ssl
        [ -f /etc/default/prometheus-node-exporter ] && mv /etc/default/prometheus-node-exporter /etc/default/prometheus-node-exporter_bak
        [ -f /etc/default/prometheus-nginx-exporter ] && mv /etc/default/prometheus-nginx-exporter /etc/default/prometheus-nginx-exporter_bak
        [ -f /etc/default/prometheus-node-exporter_pckg ] && mv /etc/default/prometheus-node-exporter_pckg /etc/default/prometheus-node-exporter
        [ -f /etc/default/prometheus-nginx-exporter_pckg ] && mv /etc/default/prometheus-nginx-exporter_pckg /etc/default/prometheus-nginx-exporter
        #[ -f /etc/nginx/sites-enabled/nginx-pckg.conf ] && ln -s /etc/nginx/sites-available/nginx-pckg.conf /etc/nginx/sites-enabled/nginx-pckg.conf
        invoke-rc.d prometheus-node-exporter stop
        invoke-rc.d prometheus-nginx-exporter stop
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.



exit 0
